# target/linux/mvebu/base-files/etc/uci-defaults/set_root_pass
#!/bin/sh

DUMMY_PASS=`tr -dc A-Za-z0-9_ < /dev/urandom | head -c 64`
echo -e "${DUMMY_PASS}\n${DUMMY_PASS}" | passwd root
